# systemd/lsyncd.service.d/override.conf
[Service]
ExecStartPre=/usr/local/bin/lsyncd-merge-includes