# debian/postinst
#!/bin/sh
set -e

# Create log directory if not exists
if [ ! -d /var/log/lsyncd ]; then
    mkdir -p /var/log/lsyncd
    chown root:adm /var/log/lsyncd
    chmod 755 /var/log/lsyncd
fi

# Create /var/log/lsyncd/lsyncd.log
if [ ! -f /var/log/lsyncd/lsyncd.log ]; then
    touch /var/log/lsyncd/lsyncd.log
    chown root:adm /var/log/lsyncd/lsyncd.log
    chmod 644 /var/log/lsyncd/lsyncd.log
fi

# Create /var/run/lsyncd.status
if [ ! -f /var/run/lsyncd.status ]; then
    touch /var/run/lsyncd.status
    chown root:root /var/run/lsyncd.status
    chmod 644 /var/run/lsyncd.status
fi

# --- Systemd integration ---
if [ -x /bin/systemctl ] || [ -x /usr/bin/systemctl ]; then
    systemctl daemon-reload || true
    systemctl try-restart lsyncd.service || true
else
    if command -v deb-systemd-helper >/dev/null 2>&1; then
        deb-systemd-helper debian-installed lsyncd.service || true
    fi
    if command -v deb-systemd-invoke >/dev/null 2>&1; then
        deb-systemd-invoke reload-or-try-restart lsyncd.service || true
    fi
fi

exit 0
